# chef/cookbooks/feedient-metrics/files/default/supervisor_metrics.conf
[program:metrics_server]
command=/usr/bin/node /opt/feedient-metrics/server.js
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=QUIT
user=www-data